# deployment/membrane.service
[Unit]
Description=Membrane Global Contextual Memory Fabric
After=network.target

[Service]
Type=simple
User=membrane
Group=membrane
WorkingDirectory=/opt/membrane
Environment=MEMBRANE_LOG_LEVEL=INFO
Environment=PYTHONPATH=/opt/membrane
ExecStart=/opt/membrane/.venv/bin/python scripts/demo_membrane.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/membrane/logs

[Install]
WantedBy=multi-user.target
